# ansible/roles/cartodb.resque/templates/cartodb_start_resque.conf.j2
# WARNING: This file is auto-generated by Ansible. Any changes will
# be overwritten the next time ansible is run.

description "Start all CartoDB resque workers"
# Copied from http://stackoverflow.com/questions/8681974/starting-multiple-upstart-instances-automatically

start on filesystem

task

env NUM_WORKERS={{ cartodb_resque_workers }}

script
  for i in `seq 1 $NUM_WORKERS`
  do
    service cartodb_resque_worker start WORKERNUM=$i
  done
end script
